# kafka/supervisor/kafka.ini
[program:kafka]
command=/opt/kafka_2.11-2.0.0/kafka.sh
autostart=true
autorestart=true
