# kafka/supervisor/zookeeper.ini
[program:zookeeper]
command=/opt/kafka_2.11-2.0.0/zoo.sh start-foreground
autostart=true
autorestart=true
